# hr.ini
[program:hr]
command=/opt/fonsview/3RD/hr_release/venv/bin/python manager.py
directory=/opt/fonsview/3RD/hr_release
autostart=true
startsecs=3
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups = 10
stdout_logfile=/opt/fonsview/3RD/hr_release/hr.log
stopasgroup=false
killasgroup=false
